# pi/salt-minion/check-salt-keys.cron
#!/bin/bash

/usr/local/bin/check-salt-keys -r 2>&1 | logger --tag check-salt-keys
