# pi/energy-savings/energy-savings.service
[Unit]
Description=Run energy savings script

[Service]
ExecStart=/usr/bin/energy-savings
Type=oneshot

[Install]
WantedBy=basic.target
