# pi/power-savings/power-savings.service
[Unit]
Description=Disable HDMI output

[Service]
ExecStart=/usr/bin/power-savings
Type=oneshot

[Install]
WantedBy=basic.target
