# pi/hdmi-off.service
[Unit]
Description=Disable HDMI output

[Service]
ExecStart=/opt/vc/bin/tvservice --off
Type=oneshot

[Install]
WantedBy=basic.target

